# sync/com.meowops.daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Meow Operations — daily digest template (macOS)

  Install:
    cp sync/com.meowops.daily-digest.plist ~/Library/LaunchAgents/ && launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.meowops.daily-digest.plist

  Unload:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.meowops.daily-digest.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.meowops.daily-digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/napster/Documents/meow-ops/sync/loop-digest.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/napster/Documents/meow-ops</string>
  <!-- launchd calendar intervals use local time; 08:30 IST is 03:00 UTC. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/napster/Library/Logs/meow-ops/daily-digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/napster/Library/Logs/meow-ops/daily-digest.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/napster</string>
  </dict>
</dict>
</plist>
